# deploy/systemd/day-trader-watchdog.timer
[Unit]
Description=Day-trader heartbeat watchdog timer (every minute during market hours)

[Timer]
# UTC times covering ET market hours under both EST (UTC-5) and EDT (UTC-4):
# Market open  09:30 ET = 13:30 UTC (EDT) / 14:30 UTC (EST)
# Market close 16:00 ET = 20:00 UTC (EDT) / 21:00 UTC (EST)
# Firing 13:00-21:30 UTC covers both DST states with margin;
# healthcheck.py gates on is_market_open() internally so over-fires
# are silent no-ops.
#
# NOTE: systemd does NOT accept full IANA timezone names like
# "America/New_York" in OnCalendar — only UTC offsets or no tz
# (= local). We use UTC explicitly via the system's default.
# Format: each field (hour, minute, second) is ranged independently.
# Cross-field ranges like `13:00..21:30` are NOT valid systemd syntax.
OnCalendar=Mon..Fri *-*-* 13..21:*:0
Persistent=false

[Install]
WantedBy=timers.target
